# PKGBUILD/steamfork-customizations/src/etc/systemd/system/systemd-logind.service.d/homedir-override.conf
# This file is needed for the hibernate serivce to be able to see the swapfile
# in /home, otherwise it will underestimate available swap space and abort

[Service]
ProtectHome=read-only
