# PKGBUILD/steamfork-device-support/src/usr/lib/udev/rules.d/99-steamfork-hdmi-watchdog.rules
KERNEL=="card*", SUBSYSTEM=="drm", ACTION=="add|change", RUN+="/usr/bin/steamfork-hdmi-watchdog"
