# PKGBUILD/steamfork-customizations-jupiter/src/usr/lib/systemd/system/steamfork-enable-swap.service
[Unit]
Description=SteamFork swap space enabler.
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/steamfork-postupdate.service
Restart=no

[Install]
WantedBy=multi-user.target
