# PKGBUILD/steamfork-device-support/src/usr/lib/systemd/user/steamfork-audio-watchdog.service
[Unit]
Description=A simple helper to route audio to hdmi or usb.
Before=wireplumber.service
ConditionUser=!root

[Service]
ExecStart=/usr/bin/steamfork-audio-watchdog
Restart=on-failure

[Install]
WantedBy=default.target
