# PKGBUILD/steamfork-device-support/src/usr/lib/systemd/user/steamfork-audio-watchdog.service
[Unit]
Description=A simple helper to route audio to hdmi or usb.
Requires=pipewire-pulse.service wireplumber.service graphical-session-pre.target
After=pipewire-pulse.service wireplumber.service graphical-session-pre.target
Before=graphical-session.target
ConditionUser=!root

[Service]
ExecStart=/usr/bin/steamfork-audio-watchdog
Restart=on-failure

[Install]
WantedBy=default.target
