# PKGBUILD/steamfork-customizations/src/etc/systemd/system/steamfork-create-overlays.service
[Unit]
Description=SteamFork Overlay - Create first boot directories
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/steamfork-create-overlays

[Install]
WantedBy=local-fs.target